# share/lxc/fuck-lxc.conf
start on started cgroup-lite
script
	echo 1 > /sys/fs/cgroup/memory/memory.use_hierarchy
end script
